# deployment/qwanyx-api.service
[Unit]
Description=QWANYX API Service
After=network.target mongodb.service
Wants=mongodb.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/qwanyx/apps/qwanyx-server/qwanyx-api
Environment="PATH=/opt/qwanyx/apps/qwanyx-server/qwanyx-api/venv/bin"
ExecStart=/opt/qwanyx/apps/qwanyx-server/qwanyx-api/venv/bin/gunicorn --workers 3 --bind 127.0.0.1:5000 app:app
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target